# pi-camera.conf
description "pi camera archiver"
start on runlevel [2]
setuid pi
setgid pi
script
  chdir /home/pi/camera/
  exec python3 archiver.py /home/pi/camera/pic 4 >> error.log 2>&1 
end script
